# q3js/rootfs/etc/services.d/q3js/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Q3JS Add-on — runs quakejs-dockerized-ssl via docker compose
# ==============================================================================

MAP_NAME=$(bashio::config 'map_name')
FRAG_LIMIT=$(bashio::config 'frag_limit')
TIME_LIMIT=$(bashio::config 'time_limit')
RCON_PASS=$(bashio::config 'rcon_password')
GAMEMODE=$(bashio::config 'gamemode')
RESTART_HOURS=$(bashio::config 'restart_hours')

WORK_DIR="/share/quakejs-dockerized-ssl"
LOG_DIR="/var/log/q3js"
GAME_LOG="${WORK_DIR}/game.log"

mkdir -p "${LOG_DIR}"

bashio::log.info "========================================="
bashio::log.info " Q3JS - Quake III Arena"
bashio::log.info " (unknown-gd/quakejs-dockerized-ssl)"
bashio::log.info "========================================="

# --------------------------------------------------------------------------
# Dynamically find the real host path for /share via docker inspect
# --------------------------------------------------------------------------
HOST_SHARE=""

CONTAINER_ID=$(grep -o '[a-f0-9]\{64\}' /proc/self/cgroup 2>/dev/null | head -1 || true)
if [ -z "${CONTAINER_ID}" ]; then
    CONTAINER_ID=$(hostname 2>/dev/null || true)
fi
bashio::log.info "Addon container ID: ${CONTAINER_ID}"

if [ -n "${CONTAINER_ID}" ]; then
    HOST_SHARE=$(docker inspect "${CONTAINER_ID}" \
        --format '{{range .Mounts}}{{if eq .Destination "/share"}}{{.Source}}{{end}}{{end}}' \
        2>/dev/null || true)
fi

if [ -z "${HOST_SHARE}" ]; then
    bashio::log.warning "Could not auto-detect host share path, falling back to /mnt/data/supervisor/share"
    HOST_SHARE="/mnt/data/supervisor/share"
fi

HOST_WORK_DIR="${HOST_SHARE}/quakejs-dockerized-ssl"
bashio::log.info "Host share path: ${HOST_SHARE}"
bashio::log.info "Host work dir:   ${HOST_WORK_DIR}"

# --------------------------------------------------------------------------
# Host LAN IP (hardcoded)
# --------------------------------------------------------------------------
HOST_IP="192.168.4.214"
bashio::log.info "Host LAN IP: ${HOST_IP}"

# --------------------------------------------------------------------------
# Stop any running containers BEFORE wiping the work dir
# --------------------------------------------------------------------------
if [ -f "${WORK_DIR}/compose.yml" ]; then
    bashio::log.info "Stopping any previous containers..."
    docker compose -f "${WORK_DIR}/compose.yml" down --remove-orphans 2>/dev/null || true
fi

# --------------------------------------------------------------------------
# Wipe and re-clone every time so Dockerfiles are always fresh for patching
# --------------------------------------------------------------------------
bashio::log.info "Wiping old clone and re-cloning fresh..."
rm -rf "${WORK_DIR}"
if ! git clone --depth 1 https://github.com/unknown-gd/quakejs-dockerized-ssl.git "${WORK_DIR}"; then
    bashio::exit.nok "git clone failed"
fi

cd "${WORK_DIR}" || bashio::exit.nok "Failed to enter work dir"

bashio::log.info "Cloned files in ${WORK_DIR}:"
ls -la "${WORK_DIR}/"

COMPOSE_FILE=$(find "${WORK_DIR}" -maxdepth 2 -name "compose.yml" | head -1 || true)
if [ -z "${COMPOSE_FILE}" ]; then
    bashio::log.error "compose.yml not found in cloned repo!"
    find "${WORK_DIR}" -maxdepth 3 | head -40
    bashio::exit.nok "compose.yml missing"
fi
bashio::log.info "Found compose file: ${COMPOSE_FILE}"

# --------------------------------------------------------------------------
# Patch compose.yml: replace ALL "./" volume paths with absolute host paths
# --------------------------------------------------------------------------
bashio::log.info "Patching compose.yml volume paths..."
sed -i "s|      - \./|      - ${HOST_WORK_DIR}/|g" "${COMPOSE_FILE}"
grep -n " - ${HOST_SHARE}" "${COMPOSE_FILE}" | head -20 || true

# --------------------------------------------------------------------------
# Add game log volume mount to dedicated server
# --------------------------------------------------------------------------
touch "${WORK_DIR}/game.log"
bashio::log.info "Patching dedicated server game log volume..."
sed -i "/dedicated-server.*server\.cfg:\/app\/base\/cpma\/server\.cfg/{
    a\      - ${HOST_WORK_DIR}/game.log:/app/game.log
}" "${COMPOSE_FILE}"

# --------------------------------------------------------------------------
# Mount pak files directly into dedicated server (bypasses DirtyInstallers bug)
# --------------------------------------------------------------------------
bashio::log.info "Mounting pak files into dedicated server..."
for PAKFILE in /share/q3js/assets/baseq3/*.pk3; do
    [ -f "${PAKFILE}" ] || continue
    PAKNAME=$(basename "${PAKFILE}")
    sed -i "/dedicated-server.*game\.log:\/app\/game\.log/{
        a\      - ${HOST_SHARE}/q3js/assets/baseq3/${PAKNAME}:/app/base/baseq3/${PAKNAME}:ro
    }" "${COMPOSE_FILE}"
    bashio::log.info "  mounted: ${PAKNAME}"
done
bashio::log.info "Dedicated server pak mounts:"
grep -A30 "q3js-dedicated-server" "${COMPOSE_FILE}" | grep "pk3" || true

# --------------------------------------------------------------------------
# Patch dedicated server: allow self-signed certs
# --------------------------------------------------------------------------
bashio::log.info "Patching NODE_TLS_REJECT_UNAUTHORIZED..."
sed -i '/q3js-dedicated-server-1:/{:loop;n;/volumes:/!b loop;i\    environment:\n      - NODE_TLS_REJECT_UNAUTHORIZED=0
}' "${COMPOSE_FILE}"

# --------------------------------------------------------------------------
# Remap port 443 -> 8443 (443 taken by Nginx Proxy Manager)
# --------------------------------------------------------------------------
bashio::log.info "Remapping port 443 -> 8443..."
sed -i 's|"443:443"|"8443:443"|g; s|- 443:443|- 8443:443|g' "${COMPOSE_FILE}"

# --------------------------------------------------------------------------
# Patch ALL Dockerfiles: strip quakejs-files from npm install
# --------------------------------------------------------------------------
bashio::log.info "Patching Dockerfiles to remove quakejs-files..."
find "${WORK_DIR}/build" -name "Dockerfile" 2>/dev/null | while read -r DF; do
    bashio::log.info "Patching: ${DF}"
    awk '{
        if (/npm install/ && !/quakejs-files/) {
            print "RUN sed -i \x27/\"quakejs-files\"/d\x27 package.json 2>/dev/null || true"
        }
        print
    }' "${DF}" > "${DF}.patched" && mv "${DF}.patched" "${DF}"
    bashio::log.info "  patch lines added: $(grep -c 'quakejs-files' "${DF}")"
done

bashio::log.info "Pruning Docker build cache..."
docker builder prune -af 2>/dev/null || true

# --------------------------------------------------------------------------
# Copy wildcard certificate from /share/q3js
# --------------------------------------------------------------------------
CERT_DIR="${WORK_DIR}/config/certificates"
mkdir -p "${CERT_DIR}"

if [ -f "/share/q3js/dual-wildcard.crt" ] && [ -f "/share/q3js/dual-wildcard.key" ]; then
    bashio::log.info "Using wildcard cert from /share/q3js/dual-wildcard.crt"
    cp "/share/q3js/dual-wildcard.crt" "${CERT_DIR}/q3js.crt"
    cp "/share/q3js/dual-wildcard.key" "${CERT_DIR}/q3js.key"
else
    bashio::log.warning "Wildcard cert not found at /share/q3js/dual-wildcard.crt — falling back to self-signed"
    openssl req -x509 -nodes -days 3650 -newkey rsa:2048 \
        -keyout "${CERT_DIR}/q3js.key" \
        -out "${CERT_DIR}/q3js.crt" \
        -subj "/CN=${HOST_IP}" 2>/dev/null
fi

cat > "${CERT_DIR}/config.json" << CERTEOF
{
    "cert": "/app/certificates/q3js.crt",
    "key": "/app/certificates/q3js.key"
}
CERTEOF

# --------------------------------------------------------------------------
# Copy pak assets into content server assets folder
# --------------------------------------------------------------------------
mkdir -p "${WORK_DIR}/assets/baseq3"
if [ -d "/share/q3js/assets/baseq3" ]; then
    bashio::log.info "Copying baseq3 pak files from /share/q3js/assets/baseq3..."
    cp -v /share/q3js/assets/baseq3/*.pk3 "${WORK_DIR}/assets/baseq3/" 2>&1 | while IFS= read -r line; do
        bashio::log.info "  ${line}"
    done
    ls -lh "${WORK_DIR}/assets/baseq3/" || true
else
    bashio::log.warning "No pak files found at /share/q3js/assets/baseq3/ — game will not load!"
fi

# --------------------------------------------------------------------------
# Write server.cfg
# --------------------------------------------------------------------------
mkdir -p "${WORK_DIR}/config/dedicated-server/baseq3"
mkdir -p "${WORK_DIR}/config/dedicated-server/cpma"

cat > "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg" << CFGEOF
seta sv_hostname "HA QuakeJS Server"
seta sv_maxclients 12
seta g_motd "Home Assistant Quake 3"
seta g_gametype 0
seta fraglimit ${FRAG_LIMIT}
seta timelimit ${TIME_LIMIT}
seta rconpassword "${RCON_PASS}"
seta g_log "/app/game.log"
seta g_logSync 1
seta bot_enable 1
CFGEOF

BOT_COUNT=$(bashio::config 'bots | length')
BOT_SKILL=$(bashio::config "bots[0].skill")
echo "seta g_spSkill ${BOT_SKILL}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"
echo "map ${MAP_NAME}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"
# addbot after map in server.cfg executes post map-load
for i in $(seq 0 $((BOT_COUNT - 1))); do
    BOT_NAME=$(bashio::config "bots[${i}].name")
    BOT_SKILL=$(bashio::config "bots[${i}].skill")
    echo "addbot ${BOT_NAME} ${BOT_SKILL}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"
done

cp "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg" \
   "${WORK_DIR}/config/dedicated-server/cpma/server.cfg"

# --------------------------------------------------------------------------
# Write .env
# --------------------------------------------------------------------------
cat > "${WORK_DIR}/.env" << ENVEOF
Q3JS_DOMAIN=${HOST_IP}
Q3JS_SERVER_GAMEMODE=${GAMEMODE}
Q3JS_SERVER_DEDICATED=1
Q3JS_SERVER_RESTART_HOURS=${RESTART_HOURS}
ENVEOF

bashio::log.info "Final volume mounts in compose.yml:"
grep -A1 "volumes:" "${COMPOSE_FILE}" | grep " - /" | head -30 || true

# --------------------------------------------------------------------------
# Build and start via docker compose
# --------------------------------------------------------------------------
bashio::log.info "Building Q3JS containers (first run takes several minutes)..."
docker compose -f "${COMPOSE_FILE}" build 2>&1 | tee "${LOG_DIR}/build.log"
if [ $? -ne 0 ]; then
    bashio::log.error "Build failed! Check /var/log/q3js/build.log"
    bashio::exit.nok "Build failed"
fi

bashio::log.info "Starting Q3JS containers..."
docker compose -f "${COMPOSE_FILE}" up -d 2>&1 | tee "${LOG_DIR}/compose.log"

# --------------------------------------------------------------------------
# Patch DirtyInstallers callback bug in dedicated server and web client JS
# --------------------------------------------------------------------------
bashio::log.info "Waiting for containers to start..."
for i in $(seq 1 30); do
    if docker inspect -f "{{.State.Running}}" q3js-dedicated-server-1 2>/dev/null | grep -q true; then
        break
    fi
    sleep 2
done

bashio::log.info "Patching DirtyInstallers in dedicated server..."
docker exec q3js-dedicated-server-1 \
    sed -i 's/return callback(/return (typeof callback === "function") \&\& callback(/g' \
    /app/build/ioq3ded_secure.js 2>/dev/null \
    && bashio::log.info "Dedicated server patch applied." \
    && docker restart q3js-dedicated-server-1 2>/dev/null \
    || bashio::log.warning "Dedicated server patch failed"

bashio::log.info "Patching DirtyInstallers in web client JS..."
Q3JS_CLIENT_JS=$(docker exec q3js-content-server-1 \
    find /app -name "__ioquake3_secure.js" 2>/dev/null | head -1 || true)
if [ -n "${Q3JS_CLIENT_JS}" ]; then
    bashio::log.info "Found client JS: ${Q3JS_CLIENT_JS}"
    docker exec q3js-content-server-1 \
        sed -i 's/return callback(/return (typeof callback === "function") \&\& callback(/g' \
        "${Q3JS_CLIENT_JS}" 2>/dev/null \
        && bashio::log.info "Web client patch applied." \
        || bashio::log.warning "Web client patch failed"
else
    bashio::log.warning "Could not find __ioquake3_secure.js in content server container"
fi

bashio::log.info "Q3JS stack is up!"
bashio::log.info "  Web client : https://${HOST_IP}:8443"
bashio::log.info "  WSS proxy  : wss://${HOST_IP}:27961"
bashio::log.info "  Content    : https://${HOST_IP}:9000"

# --------------------------------------------------------------------------
# Pipe docker logs into game.log (g_log is relative inside container,
# so we capture stdout which contains all game events)
# --------------------------------------------------------------------------
bashio::log.info "Starting docker log -> game.log pipe..."
(
    for i in $(seq 1 60); do
        if docker inspect -f "{{.State.Running}}" q3js-dedicated-server-1 2>/dev/null | grep -q true; then
            break
        fi
        sleep 2
    done
    bashio::log.info "Piping q3js-dedicated-server-1 logs to ${GAME_LOG}"
    docker logs --since 1s -f q3js-dedicated-server-1 2>&1 | while IFS= read -r line; do
        echo "${line}" >> "${GAME_LOG}"
        bashio::log.info "[q3ds] ${line}"
    done
) &

# --------------------------------------------------------------------------
# Start metrics server
# --------------------------------------------------------------------------
bashio::log.info "Starting metrics server on :8090..."
Q3_LOG_FILE="${GAME_LOG}" \
METRICS_PORT=8090 \
RCON_PASS="${RCON_PASS}" \
    /usr/bin/q3js-metrics > "${LOG_DIR}/metrics.log" 2>&1 &

# Keep alive
while true; do
    sleep 30
done
