# q3js/rootfs/etc/services.d/q3js/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Q3JS Add-on — runs quakejs-dockerized-ssl via docker compose
# ==============================================================================

MAP_NAME=$(bashio::config 'map_name')
FRAG_LIMIT=$(bashio::config 'frag_limit')
TIME_LIMIT=$(bashio::config 'time_limit')
RCON_PASS=$(bashio::config 'rcon_password')
GAMEMODE=$(bashio::config 'gamemode')
RESTART_HOURS=$(bashio::config 'restart_hours')

WORK_DIR="/data/quakejs-dockerized-ssl"
LOG_DIR="/var/log/q3js"
GAME_LOG="${LOG_DIR}/game.log"

mkdir -p "${LOG_DIR}"

bashio::log.info "========================================="
bashio::log.info " Q3JS - Quake III Arena"
bashio::log.info " (unknown-gd/quakejs-dockerized-ssl)"
bashio::log.info "========================================="

# Clone or update
if [ ! -d "${WORK_DIR}/.git" ]; then
    bashio::log.info "Cloning quakejs-dockerized-ssl..."
    git clone --depth 1 https://github.com/unknown-gd/quakejs-dockerized-ssl.git "${WORK_DIR}"
else
    bashio::log.info "Repo already cloned, skipping update."
fi

cd "${WORK_DIR}" || bashio::exit.nok "Failed to enter work dir"

# --------------------------------------------------------------------------
# Patch: quakejs-files was removed from npm. Each service Dockerfile clones
# https://github.com/unknown-gd/quakejs inside the build container, so we
# cannot patch package.json on the host. Instead we inject a sed RUN step
# directly into each Dockerfile before npm install runs.
# We also prune the Docker build cache so the patched Dockerfile is used.
# --------------------------------------------------------------------------
bashio::log.info "Patching Dockerfiles: removing quakejs-files from npm install..."
find "${WORK_DIR}/build" -name "Dockerfile" 2>/dev/null | while read -r DF; do
    if grep -q "npm install" "${DF}" && ! grep -q "quakejs-files" "${DF}"; then
        # Match ANY line containing npm install and prepend the sed removal
        sed -i 's|RUN \(.*\)npm install|RUN sed -i /"quakejs-files"/d package.json \&\& \1npm install|g' "${DF}"
        bashio::log.info "Patched ${DF}"
    fi
done

bashio::log.info "Pruning Docker build cache to force fresh build..."
docker builder prune -af 2>/dev/null || true

# --------------------------------------------------------------------------
# Generate self-signed certificate
# --------------------------------------------------------------------------
CERT_DIR="${WORK_DIR}/config/certificates"
mkdir -p "${CERT_DIR}"

if [ ! -f "${CERT_DIR}/q3js.crt" ]; then
    bashio::log.info "Generating self-signed certificate..."
    openssl req -x509 -nodes -days 3650 -newkey rsa:2048 \
        -keyout "${CERT_DIR}/q3js.key" \
        -out "${CERT_DIR}/q3js.crt" \
        -subj "/CN=homeassistant.local" 2>/dev/null
fi

cat > "${CERT_DIR}/config.json" << CERTEOF
{
    "cert": "/app/certificates/q3js.crt",
    "key": "/app/certificates/q3js.key"
}
CERTEOF

# --------------------------------------------------------------------------
# Copy pak assets from /share
# --------------------------------------------------------------------------
mkdir -p "${WORK_DIR}/assets"
if [ -d "/share/q3js/assets" ]; then
    bashio::log.info "Copying assets from /share/q3js/assets..."
    cp -rn /share/q3js/assets/. "${WORK_DIR}/assets/" 2>/dev/null || true
fi

# --------------------------------------------------------------------------
# Write server.cfg
# --------------------------------------------------------------------------
mkdir -p "${WORK_DIR}/config/dedicated-server/baseq3"
mkdir -p "${WORK_DIR}/config/dedicated-server/cpma"

cat > "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg" << CFGEOF
seta sv_hostname "HA QuakeJS Server"
seta sv_maxclients 12
seta g_motd "Home Assistant Quake 3"
seta g_gametype 0
seta fraglimit ${FRAG_LIMIT}
seta timelimit ${TIME_LIMIT}
seta rconpassword "${RCON_PASS}"
seta g_log "${GAME_LOG}"
seta g_logSync 1
seta bot_enable 1
CFGEOF

BOT_COUNT=$(bashio::config 'bots | length')
for i in $(seq 0 $((BOT_COUNT - 1))); do
    BOT_NAME=$(bashio::config "bots[${i}].name")
    BOT_SKILL=$(bashio::config "bots[${i}].skill")
    echo "addbot ${BOT_NAME} ${BOT_SKILL}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"
done
echo "map ${MAP_NAME}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"

cp "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg" \
   "${WORK_DIR}/config/dedicated-server/cpma/server.cfg"

# --------------------------------------------------------------------------
# Write .env
# --------------------------------------------------------------------------
cat > "${WORK_DIR}/.env" << ENVEOF
Q3JS_DOMAIN=homeassistant.local
Q3JS_SERVER_GAMEMODE=${GAMEMODE}
Q3JS_SERVER_DEDICATED=1
Q3JS_SERVER_RESTART_HOURS=${RESTART_HOURS}
ENVEOF

# --------------------------------------------------------------------------
# Build and start via docker compose
# --------------------------------------------------------------------------
bashio::log.info "Stopping any previous containers..."
docker compose down --remove-orphans 2>/dev/null || true

bashio::log.info "Building Q3JS containers (first run takes several minutes)..."
docker compose build 2>&1 | tee "${LOG_DIR}/build.log"
if [ $? -ne 0 ]; then
    bashio::log.error "docker compose build failed! Check /var/log/q3js/build.log"
    bashio::exit.nok "Build failed"
fi

bashio::log.info "Starting Q3JS containers..."
docker compose up -d 2>&1 | tee "${LOG_DIR}/compose.log"

bashio::log.info "Q3JS stack is up!"
bashio::log.info "  Web client : https://homeassistant.local:443"
bashio::log.info "  WSS proxy  : wss://homeassistant.local:27961"
bashio::log.info "  Content    : https://homeassistant.local:9000"

# --------------------------------------------------------------------------
# Start metrics server
# --------------------------------------------------------------------------
bashio::log.info "Starting metrics server on :8090..."
Q3_LOG_FILE="${GAME_LOG}" \
METRICS_PORT=8090 \
RCON_PASS="${RCON_PASS}" \
    /usr/bin/q3js-metrics > "${LOG_DIR}/metrics.log" 2>&1 &

# Keep alive
while true; do
    sleep 30
done
