# q3js/rootfs/etc/services.d/q3js/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Q3JS Add-on — runs quakejs-dockerized-ssl via docker compose
# ==============================================================================

MAP_NAME=$(bashio::config 'map_name')
FRAG_LIMIT=$(bashio::config 'frag_limit')
TIME_LIMIT=$(bashio::config 'time_limit')
RCON_PASS=$(bashio::config 'rcon_password')
GAMEMODE=$(bashio::config 'gamemode')
RESTART_HOURS=$(bashio::config 'restart_hours')
CONNECT_SERVER=$(bashio::config 'connect_server' 2>/dev/null || echo "")

WORK_DIR="/share/quakejs-dockerized-ssl"
LOG_DIR="/var/log/q3js"
GAME_LOG="${WORK_DIR}/game.log"

mkdir -p "${LOG_DIR}"

bashio::log.info "========================================="
bashio::log.info " Q3JS - Quake III Arena"
bashio::log.info " (rogerquake/quakejs-dockerized-ssl)"
bashio::log.info "========================================="

# --------------------------------------------------------------------------
# Dynamically find the real host path for /share via docker inspect
# --------------------------------------------------------------------------
HOST_SHARE=""

CONTAINER_ID=$(grep -o '[a-f0-9]\{64\}' /proc/self/cgroup 2>/dev/null | head -1 || true)
if [ -z "${CONTAINER_ID}" ]; then
    CONTAINER_ID=$(hostname 2>/dev/null || true)
fi
bashio::log.info "Addon container ID: ${CONTAINER_ID}"

if [ -n "${CONTAINER_ID}" ]; then
    HOST_SHARE=$(docker inspect "${CONTAINER_ID}" \
        --format '{{range .Mounts}}{{if eq .Destination "/share"}}{{.Source}}{{end}}{{end}}' \
        2>/dev/null || true)
fi

if [ -z "${HOST_SHARE}" ]; then
    bashio::log.warning "Could not auto-detect host share path, falling back to /mnt/data/supervisor/share"
    HOST_SHARE="/mnt/data/supervisor/share"
fi

HOST_WORK_DIR="${HOST_SHARE}/quakejs-dockerized-ssl"
bashio::log.info "Host share path: ${HOST_SHARE}"
bashio::log.info "Host work dir:   ${HOST_WORK_DIR}"

# --------------------------------------------------------------------------
# Host LAN IP (hardcoded)
# --------------------------------------------------------------------------
HOST_IP="homeassistant-test.roger.lan"
bashio::log.info "Host LAN IP: ${HOST_IP}"

# --------------------------------------------------------------------------
# Stop any running containers BEFORE wiping the work dir
# --------------------------------------------------------------------------
if [ -f "${WORK_DIR}/compose.yml" ]; then
    bashio::log.info "Stopping any previous containers..."
    docker compose -f "${WORK_DIR}/compose.yml" down --remove-orphans 2>/dev/null || true
fi

# --------------------------------------------------------------------------
# Clone fork — DirtyInstallers fix and username prompt removal
# are already baked into the fork's Dockerfiles and index.ejs
# --------------------------------------------------------------------------
bashio::log.info "Cloning fork..."
rm -rf "${WORK_DIR}"
if ! git clone --depth 1 https://github.com/rogerquake/quakejs-dockerized-ssl.git "${WORK_DIR}"; then
    bashio::exit.nok "git clone failed"
fi

cd "${WORK_DIR}" || bashio::exit.nok "Failed to enter work dir"

bashio::log.info "Cloned files in ${WORK_DIR}:"
ls -la "${WORK_DIR}/"

COMPOSE_FILE=$(find "${WORK_DIR}" -maxdepth 2 -name "compose.yml" | head -1 || true)
if [ -z "${COMPOSE_FILE}" ]; then
    bashio::log.error "compose.yml not found in cloned repo!"
    find "${WORK_DIR}" -maxdepth 3 | head -40
    bashio::exit.nok "compose.yml missing"
fi
bashio::log.info "Found compose file: ${COMPOSE_FILE}"

# --------------------------------------------------------------------------
# Patch compose.yml: replace all "./" volume paths with absolute host paths
# --------------------------------------------------------------------------
bashio::log.info "Patching compose.yml volume paths..."
sed -i "s|      - \./|      - ${HOST_WORK_DIR}/|g" "${COMPOSE_FILE}"

# --------------------------------------------------------------------------
# Add game log + pak mounts to dedicated server using awk
# --------------------------------------------------------------------------
touch "${WORK_DIR}/game.log"
bashio::log.info "Adding volume mounts to dedicated server..."

# Build a temp file containing all extra volume lines to inject
MOUNTS_FILE="${WORK_DIR}/.extra_mounts"
rm -f "${MOUNTS_FILE}"

echo "      - ${HOST_WORK_DIR}/game.log:/app/game.log" >> "${MOUNTS_FILE}"

for PAKFILE in /share/q3js/assets/baseq3/*.pk3; do
    [ -f "${PAKFILE}" ] || continue
    PAKNAME=$(basename "${PAKFILE}")
    echo "      - ${HOST_SHARE}/q3js/assets/baseq3/${PAKNAME}:/app/base/baseq3/${PAKNAME}:ro" >> "${MOUNTS_FILE}"
done

for MOD_DIR in /share/q3js/assets/*/; do
    MOD_NAME=$(basename "${MOD_DIR}")
    [ "${MOD_NAME}" = "baseq3" ] && continue
    for PAKFILE in "${MOD_DIR}"*.pk3; do
        [ -f "${PAKFILE}" ] || continue
        PAKNAME=$(basename "${PAKFILE}")
        echo "      - ${HOST_SHARE}/q3js/assets/${MOD_NAME}/${PAKNAME}:/app/base/${MOD_NAME}/${PAKNAME}:ro" >> "${MOUNTS_FILE}"
    done
done

bashio::log.info "Mounts to inject: $(wc -l < "${MOUNTS_FILE}")"

# Inject mounts after the last existing volume line in the dedicated server block
# awk tracks which service we are in and where the last volume line is
awk -v mfile="${MOUNTS_FILE}" '
BEGIN { while ((getline ml < mfile) > 0) mounts[++nm] = ml }
/q3js-dedicated-server-1:/ { in_ds=1 }
in_ds && /^  [a-zA-Z]/ && !/q3js-dedicated-server-1/ { in_ds=0; in_vol=0 }
in_ds && /^    volumes:/ { in_vol=1 }
{
    buf[NR] = $0
    if (in_vol && /^      - /) last_vol = NR
}
END {
    for (i=1; i<=NR; i++) {
        print buf[i]
        if (i == last_vol) {
            for (j=1; j<=nm; j++) print mounts[j]
        }
    }
}
' "${COMPOSE_FILE}" > "${COMPOSE_FILE}.tmp" && mv "${COMPOSE_FILE}.tmp" "${COMPOSE_FILE}"

bashio::log.info "Volume mounts injected"
grep -A60 "q3js-dedicated-server-1:" "${COMPOSE_FILE}" | grep "pk3\|game.log" | head -20 || true


# NODE_TLS_REJECT_UNAUTHORIZED and port 8443 are now set directly in the
# fork's compose.yml — no runtime patching needed

# --------------------------------------------------------------------------
# Patch index.ejs: inject connect server if configured
# --------------------------------------------------------------------------
INDEX_EJS=$(find "${WORK_DIR}" -name "index.ejs" | head -1)
if [ -n "${INDEX_EJS}" ]; then
    if [ -n "${CONNECT_SERVER}" ]; then
        bashio::log.info "Patching index.ejs for auto-connect: ${CONNECT_SERVER}"
        sed -i "s|// args.push.apply.*+connect.*|args.push.apply( args, [ \"+connect\", \"${CONNECT_SERVER}\" ] );|" "${INDEX_EJS}"
        bashio::log.info "index.ejs connect patch applied"
    else
        bashio::log.warning "connect_server not set — skipping auto-connect patch"
    fi
else
    bashio::log.warning "index.ejs not found"
fi


# --------------------------------------------------------------------------
# Copy wildcard certificate from /share/q3js
# --------------------------------------------------------------------------
CERT_DIR="${WORK_DIR}/config/certificates"
mkdir -p "${CERT_DIR}"

if [ -f "/share/q3js/dual-wildcard.crt" ] && [ -f "/share/q3js/dual-wildcard.key" ]; then
    bashio::log.info "Using wildcard cert from /share/q3js/dual-wildcard.crt"
    cp "/share/q3js/dual-wildcard.crt" "${CERT_DIR}/q3js.crt"
    cp "/share/q3js/dual-wildcard.key" "${CERT_DIR}/q3js.key"
else
    bashio::log.warning "Wildcard cert not found — falling back to self-signed"
    openssl req -x509 -nodes -days 3650 -newkey rsa:2048 \
        -keyout "${CERT_DIR}/q3js.key" \
        -out "${CERT_DIR}/q3js.crt" \
        -subj "/CN=${HOST_IP}" 2>/dev/null
fi

cat > "${CERT_DIR}/config.json" << CERTEOF
{
    "cert": "/app/certificates/q3js.crt",
    "key": "/app/certificates/q3js.key"
}
CERTEOF

# --------------------------------------------------------------------------
# Copy pak assets into content server assets folder
# --------------------------------------------------------------------------
mkdir -p "${WORK_DIR}/assets/baseq3"
if [ -d "/share/q3js/assets/baseq3" ]; then
    bashio::log.info "Copying baseq3 pak files..."
    cp -v /share/q3js/assets/baseq3/*.pk3 "${WORK_DIR}/assets/baseq3/" 2>&1 | while IFS= read -r line; do
        bashio::log.info "  ${line}"
    done
else
    bashio::log.warning "No pak files found at /share/q3js/assets/baseq3/ — game will not load!"
fi

# Copy all mod asset folders generically (ra3, cpma, osp, and any future mods)
for MOD_DIR in /share/q3js/assets/*/; do
    MOD_NAME=$(basename "${MOD_DIR}")
    [ "${MOD_NAME}" = "baseq3" ] && continue  # already copied above
    if ls "${MOD_DIR}"*.pk3 >/dev/null 2>&1; then
        bashio::log.info "Copying ${MOD_NAME} pak files to content server..."
        mkdir -p "${WORK_DIR}/assets/${MOD_NAME}"
        cp "${MOD_DIR}"*.pk3 "${WORK_DIR}/assets/${MOD_NAME}/"
    fi
done

# --------------------------------------------------------------------------
# Write server.cfg
# --------------------------------------------------------------------------
mkdir -p "${WORK_DIR}/config/dedicated-server/baseq3"
mkdir -p "${WORK_DIR}/config/dedicated-server/cpma"
mkdir -p "${WORK_DIR}/config/dedicated-server/ra3"

BOT_COUNT=$(bashio::config 'bots | length')
BOT_SKILL=$(bashio::config "bots[0].skill")

if [ "${GAMEMODE}" = "cpma" ]; then
    cat > "${WORK_DIR}/config/dedicated-server/cpma/server.cfg" << CPMACFG
seta sv_hostname "HA QuakeJS CPMA Server"
seta sv_maxclients 12
seta g_motd "Home Assistant CPMA"
seta fs_game "cpma"
seta g_gametype 1
seta fraglimit ${FRAG_LIMIT}
seta timelimit ${TIME_LIMIT}
seta rconpassword "${RCON_PASS}"
seta g_log "/app/game.log"
seta g_logSync 1
seta bot_enable 1
CPMACFG
    echo "seta g_spSkill ${BOT_SKILL}" >> "${WORK_DIR}/config/dedicated-server/cpma/server.cfg"
    echo "map ${MAP_NAME}" >> "${WORK_DIR}/config/dedicated-server/cpma/server.cfg"
    for i in $(seq 0 $((BOT_COUNT - 1))); do
        BOT_NAME=$(bashio::config "bots[${i}].name")
        BOT_SKILL=$(bashio::config "bots[${i}].skill")
        echo "addbot ${BOT_NAME} ${BOT_SKILL}" >> "${WORK_DIR}/config/dedicated-server/cpma/server.cfg"
    done
    cp "${WORK_DIR}/config/dedicated-server/cpma/server.cfg" \
       "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"

elif [ "${GAMEMODE}" = "ra3" ]; then
    cat > "${WORK_DIR}/config/dedicated-server/ra3/server.cfg" << RA3CFG
seta sv_hostname "HA QuakeJS RA3 Server"
seta sv_maxclients 12
seta g_motd "Home Assistant RA3"
seta fs_game "ra3"
seta g_gametype 1
seta fraglimit ${FRAG_LIMIT}
seta timelimit ${TIME_LIMIT}
seta rconpassword "${RCON_PASS}"
seta g_log "/app/game.log"
seta g_logSync 1
RA3CFG
    echo "map ${MAP_NAME}" >> "${WORK_DIR}/config/dedicated-server/ra3/server.cfg"
    cp "${WORK_DIR}/config/dedicated-server/ra3/server.cfg" \
       "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"
    cp "${WORK_DIR}/config/dedicated-server/ra3/server.cfg" \
       "${WORK_DIR}/config/dedicated-server/cpma/server.cfg"

else
    cat > "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg" << CFGEOF
seta sv_hostname "HA QuakeJS Server"
seta sv_maxclients 12
seta g_motd "Home Assistant Quake 3"
seta g_gametype 0
seta fraglimit ${FRAG_LIMIT}
seta timelimit ${TIME_LIMIT}
seta rconpassword "${RCON_PASS}"
seta g_log "/app/game.log"
seta g_logSync 1
seta bot_enable 1
CFGEOF
    echo "seta g_spSkill ${BOT_SKILL}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"
    echo "map ${MAP_NAME}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"
    for i in $(seq 0 $((BOT_COUNT - 1))); do
        BOT_NAME=$(bashio::config "bots[${i}].name")
        BOT_SKILL=$(bashio::config "bots[${i}].skill")
        echo "addbot ${BOT_NAME} ${BOT_SKILL}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"
    done
    cp "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg" \
       "${WORK_DIR}/config/dedicated-server/cpma/server.cfg"
fi

# --------------------------------------------------------------------------
# Write .env
# --------------------------------------------------------------------------
cat > "${WORK_DIR}/.env" << ENVEOF
Q3JS_DOMAIN=${HOST_IP}
Q3JS_SERVER_GAMEMODE=${GAMEMODE}
Q3JS_SERVER_DEDICATED=1
Q3JS_SERVER_RESTART_HOURS=${RESTART_HOURS}
ENVEOF

bashio::log.info "Final volume mounts in compose.yml:"
grep -A1 "volumes:" "${COMPOSE_FILE}" | grep " - /" | head -30 || true

# --------------------------------------------------------------------------
# Build and start via docker compose
# --------------------------------------------------------------------------
bashio::log.info "Building Q3JS containers (first run takes several minutes)..."
docker compose -f "${COMPOSE_FILE}" build 2>&1 | tee "${LOG_DIR}/build.log"
if [ $? -ne 0 ]; then
    bashio::log.error "Build failed! Check /var/log/q3js/build.log"
    bashio::exit.nok "Build failed"
fi

bashio::log.info "Starting Q3JS containers..."
docker compose -f "${COMPOSE_FILE}" up -d 2>&1 | tee "${LOG_DIR}/compose.log"

# --------------------------------------------------------------------------
bashio::log.info "Q3JS stack is up!"
bashio::log.info "  Web client : https://${HOST_IP}:8443"
bashio::log.info "  WSS proxy  : wss://${HOST_IP}:27961"
bashio::log.info "  Content    : https://${HOST_IP}:9000"

# --------------------------------------------------------------------------
# Pipe docker logs into game.log
# --------------------------------------------------------------------------
bashio::log.info "Starting docker log -> game.log pipe..."
(
    for i in $(seq 1 60); do
        if docker inspect -f "{{.State.Running}}" q3js-dedicated-server-1 2>/dev/null | grep -q true; then
            break
        fi
        sleep 2
    done
    bashio::log.info "Piping q3js-dedicated-server-1 logs to ${GAME_LOG}"
    docker logs --since 1s -f q3js-dedicated-server-1 2>&1 | while IFS= read -r line; do
        echo "${line}" >> "${GAME_LOG}"
        bashio::log.info "[q3ds] ${line}"
    done
) &

# --------------------------------------------------------------------------
# Start metrics server
# --------------------------------------------------------------------------
bashio::log.info "Starting metrics server on :8090..."
Q3_LOG_FILE="${GAME_LOG}" \
METRICS_PORT=8090 \
RCON_PASS="${RCON_PASS}" \
    /usr/bin/q3js-metrics > "${LOG_DIR}/metrics.log" 2>&1 &

# Keep alive
while true; do
    sleep 30
done
