# q3js/rootfs/etc/services.d/q3js/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Q3JS Add-on — runs quakejs-dockerized-ssl via docker compose
# ==============================================================================

MAP_NAME=$(bashio::config 'map_name')
FRAG_LIMIT=$(bashio::config 'frag_limit')
TIME_LIMIT=$(bashio::config 'time_limit')
RCON_PASS=$(bashio::config 'rcon_password')
GAMEMODE=$(bashio::config 'gamemode')
RESTART_HOURS=$(bashio::config 'restart_hours')

WORK_DIR="/share/quakejs-dockerized-ssl"
LOG_DIR="/var/log/q3js"
GAME_LOG="${LOG_DIR}/game.log"

mkdir -p "${LOG_DIR}"

bashio::log.info "========================================="
bashio::log.info " Q3JS - Quake III Arena"
bashio::log.info " (unknown-gd/quakejs-dockerized-ssl)"
bashio::log.info "========================================="

# --------------------------------------------------------------------------
# Dynamically find the real host path for /share via docker inspect
# --------------------------------------------------------------------------
HOST_SHARE=""

CONTAINER_ID=$(grep -o '[a-f0-9]\{64\}' /proc/self/cgroup 2>/dev/null | head -1 || true)
if [ -z "${CONTAINER_ID}" ]; then
    CONTAINER_ID=$(hostname 2>/dev/null || true)
fi
bashio::log.info "Addon container ID: ${CONTAINER_ID}"

if [ -n "${CONTAINER_ID}" ]; then
    HOST_SHARE=$(docker inspect "${CONTAINER_ID}" \
        --format '{{range .Mounts}}{{if eq .Destination "/share"}}{{.Source}}{{end}}{{end}}' \
        2>/dev/null || true)
fi

if [ -z "${HOST_SHARE}" ]; then
    bashio::log.warning "Could not auto-detect host share path, falling back to /mnt/data/supervisor/share"
    HOST_SHARE="/mnt/data/supervisor/share"
fi

HOST_WORK_DIR="${HOST_SHARE}/quakejs-dockerized-ssl"
bashio::log.info "Host share path: ${HOST_SHARE}"
bashio::log.info "Host work dir:   ${HOST_WORK_DIR}"

# --------------------------------------------------------------------------
# Stop any running containers BEFORE wiping the work dir
# --------------------------------------------------------------------------
if [ -f "${WORK_DIR}/docker-compose.yml" ]; then
    bashio::log.info "Stopping any previous containers..."
    docker compose -f "${WORK_DIR}/docker-compose.yml" down --remove-orphans 2>/dev/null || true
fi

# --------------------------------------------------------------------------
# Wipe and re-clone every time so Dockerfiles are always fresh for patching
# --------------------------------------------------------------------------
bashio::log.info "Wiping old clone and re-cloning fresh..."
rm -rf "${WORK_DIR}"
git clone --depth 1 https://github.com/unknown-gd/quakejs-dockerized-ssl.git "${WORK_DIR}"

cd "${WORK_DIR}" || bashio::exit.nok "Failed to enter work dir"

# --------------------------------------------------------------------------
# Patch docker-compose.yml: replace ALL "./" volume paths with absolute
# host paths so the host Docker daemon can resolve them correctly.
# (Without this, host Docker resolves "./" as "/share/..." which is read-only)
# --------------------------------------------------------------------------
bashio::log.info "Patching docker-compose.yml volume paths..."
bashio::log.info "  Replacing './' with '${HOST_WORK_DIR}/'"
sed -i "s|      - \./|      - ${HOST_WORK_DIR}/|g" "${WORK_DIR}/docker-compose.yml"
bashio::log.info "Patched volume lines:"
grep -n "HOST_WORK_DIR\|${HOST_WORK_DIR}" "${WORK_DIR}/docker-compose.yml" || \
    grep -n " - ${HOST_SHARE}" "${WORK_DIR}/docker-compose.yml" || true

# --------------------------------------------------------------------------
# Patch ALL Dockerfiles: insert sed before npm install to strip quakejs-files
# --------------------------------------------------------------------------
bashio::log.info "Patching Dockerfiles to remove quakejs-files..."
find "${WORK_DIR}/build" -name "Dockerfile" 2>/dev/null | while read -r DF; do
    bashio::log.info "Patching: ${DF}"
    awk '{
        if (/npm install/ && !/quakejs-files/) {
            print "RUN sed -i \x27/\"quakejs-files\"/d\x27 package.json 2>/dev/null || true"
        }
        print
    }' "${DF}" > "${DF}.patched" && mv "${DF}.patched" "${DF}"
    bashio::log.info "  patch lines added: $(grep -c 'quakejs-files' "${DF}")"
done

bashio::log.info "Pruning Docker build cache..."
docker builder prune -af 2>/dev/null || true

# --------------------------------------------------------------------------
# Generate self-signed certificate
# --------------------------------------------------------------------------
CERT_DIR="${WORK_DIR}/config/certificates"
mkdir -p "${CERT_DIR}"

if [ ! -f "${CERT_DIR}/q3js.crt" ]; then
    bashio::log.info "Generating self-signed certificate..."
    openssl req -x509 -nodes -days 3650 -newkey rsa:2048 \
        -keyout "${CERT_DIR}/q3js.key" \
        -out "${CERT_DIR}/q3js.crt" \
        -subj "/CN=homeassistant.local" 2>/dev/null
fi

cat > "${CERT_DIR}/config.json" << CERTEOF
{
    "cert": "/app/certificates/q3js.crt",
    "key": "/app/certificates/q3js.key"
}
CERTEOF

# --------------------------------------------------------------------------
# Copy pak assets from /share/q3js/assets if present
# --------------------------------------------------------------------------
mkdir -p "${WORK_DIR}/assets"
if [ -d "/share/q3js/assets" ]; then
    bashio::log.info "Copying assets from /share/q3js/assets..."
    cp -rn /share/q3js/assets/. "${WORK_DIR}/assets/" 2>/dev/null || true
fi

# --------------------------------------------------------------------------
# Write server.cfg
# --------------------------------------------------------------------------
mkdir -p "${WORK_DIR}/config/dedicated-server/baseq3"
mkdir -p "${WORK_DIR}/config/dedicated-server/cpma"

cat > "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg" << CFGEOF
seta sv_hostname "HA QuakeJS Server"
seta sv_maxclients 12
seta g_motd "Home Assistant Quake 3"
seta g_gametype 0
seta fraglimit ${FRAG_LIMIT}
seta timelimit ${TIME_LIMIT}
seta rconpassword "${RCON_PASS}"
seta g_log "${GAME_LOG}"
seta g_logSync 1
seta bot_enable 1
CFGEOF

BOT_COUNT=$(bashio::config 'bots | length')
for i in $(seq 0 $((BOT_COUNT - 1))); do
    BOT_NAME=$(bashio::config "bots[${i}].name")
    BOT_SKILL=$(bashio::config "bots[${i}].skill")
    echo "addbot ${BOT_NAME} ${BOT_SKILL}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"
done
echo "map ${MAP_NAME}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"

cp "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg" \
   "${WORK_DIR}/config/dedicated-server/cpma/server.cfg"

# --------------------------------------------------------------------------
# Write .env
# --------------------------------------------------------------------------
cat > "${WORK_DIR}/.env" << ENVEOF
Q3JS_DOMAIN=homeassistant.local
Q3JS_SERVER_GAMEMODE=${GAMEMODE}
Q3JS_SERVER_DEDICATED=1
Q3JS_SERVER_RESTART_HOURS=${RESTART_HOURS}
ENVEOF

# --------------------------------------------------------------------------
# Log the final docker-compose.yml volume section for debugging
# --------------------------------------------------------------------------
bashio::log.info "Final volume mounts in docker-compose.yml:"
grep -A1 "volumes:" "${WORK_DIR}/docker-compose.yml" | grep " - /" | head -20 || true

# --------------------------------------------------------------------------
# Build and start via docker compose
# --------------------------------------------------------------------------
bashio::log.info "Building Q3JS containers (first run takes several minutes)..."
docker compose build 2>&1 | tee "${LOG_DIR}/build.log"
if [ $? -ne 0 ]; then
    bashio::log.error "Build failed! Check /var/log/q3js/build.log"
    bashio::exit.nok "Build failed"
fi

bashio::log.info "Starting Q3JS containers..."
docker compose up -d 2>&1 | tee "${LOG_DIR}/compose.log"

bashio::log.info "Q3JS stack is up!"
bashio::log.info "  Web client : https://homeassistant.local:443"
bashio::log.info "  WSS proxy  : wss://homeassistant.local:27961"
bashio::log.info "  Content    : https://homeassistant.local:9000"

# --------------------------------------------------------------------------
# Start metrics server
# --------------------------------------------------------------------------
bashio::log.info "Starting metrics server on :8090..."
Q3_LOG_FILE="${GAME_LOG}" \
METRICS_PORT=8090 \
RCON_PASS="${RCON_PASS}" \
    /usr/bin/q3js-metrics > "${LOG_DIR}/metrics.log" 2>&1 &

# Keep alive
while true; do
    sleep 30
done
