# q3js/rootfs/etc/services.d/q3js/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Q3JS Add-on — runs quakejs-dockerized-ssl via docker compose
# ==============================================================================

MAP_NAME=$(bashio::config 'map_name')
FRAG_LIMIT=$(bashio::config 'frag_limit')
TIME_LIMIT=$(bashio::config 'time_limit')
RCON_PASS=$(bashio::config 'rcon_password')
GAMEMODE=$(bashio::config 'gamemode')
RESTART_HOURS=$(bashio::config 'restart_hours')

WORK_DIR="/share/quakejs-dockerized-ssl"
LOG_DIR="/var/log/q3js"
GAME_LOG="${WORK_DIR}/game.log"

mkdir -p "${LOG_DIR}"

bashio::log.info "========================================="
bashio::log.info " Q3JS - Quake III Arena"
bashio::log.info " (unknown-gd/quakejs-dockerized-ssl)"
bashio::log.info "========================================="

# --------------------------------------------------------------------------
# Dynamically find the real host path for /share via docker inspect
# --------------------------------------------------------------------------
HOST_SHARE=""

CONTAINER_ID=$(grep -o '[a-f0-9]\{64\}' /proc/self/cgroup 2>/dev/null | head -1 || true)
if [ -z "${CONTAINER_ID}" ]; then
    CONTAINER_ID=$(hostname 2>/dev/null || true)
fi
bashio::log.info "Addon container ID: ${CONTAINER_ID}"

if [ -n "${CONTAINER_ID}" ]; then
    HOST_SHARE=$(docker inspect "${CONTAINER_ID}" \
        --format '{{range .Mounts}}{{if eq .Destination "/share"}}{{.Source}}{{end}}{{end}}' \
        2>/dev/null || true)
fi

if [ -z "${HOST_SHARE}" ]; then
    bashio::log.warning "Could not auto-detect host share path, falling back to /mnt/data/supervisor/share"
    HOST_SHARE="/mnt/data/supervisor/share"
fi

HOST_WORK_DIR="${HOST_SHARE}/quakejs-dockerized-ssl"
bashio::log.info "Host share path: ${HOST_SHARE}"
bashio::log.info "Host work dir:   ${HOST_WORK_DIR}"

# --------------------------------------------------------------------------
# Detect host LAN IP via Supervisor API (avoids Docker network IPs)
# --------------------------------------------------------------------------
HOST_IP=$(curl -sf \
    -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
    http://supervisor/network/interface/default/info \
    | python3 -c "
import json,sys
d=json.load(sys.stdin)
addrs=d.get('data',{}).get('ipv4',{}).get('address',[])
if addrs: print(addrs[0].split('/')[0])
" 2>/dev/null || true)
if [ -z "${HOST_IP}" ]; then
    bashio::log.warning "Supervisor API IP detection failed, trying ip route..."
    HOST_IP=$(ip route get 8.8.8.8 2>/dev/null | awk '/via/{for(i=1;i<=NF;i++) if ($i=="src") print $(i+1)}' | head -1 || true)
fi
if [ -z "${HOST_IP}" ]; then
    HOST_IP="homeassistant.local"
fi
bashio::log.info "Detected host LAN IP: ${HOST_IP}"

# --------------------------------------------------------------------------
# Stop any running containers BEFORE wiping the work dir
# --------------------------------------------------------------------------
if [ -f "${WORK_DIR}/compose.yml" ]; then
    bashio::log.info "Stopping any previous containers..."
    docker compose -f "${WORK_DIR}/compose.yml" down --remove-orphans 2>/dev/null || true
fi

# --------------------------------------------------------------------------
# Wipe and re-clone every time so Dockerfiles are always fresh for patching
# --------------------------------------------------------------------------
bashio::log.info "Wiping old clone and re-cloning fresh..."
rm -rf "${WORK_DIR}"
if ! git clone --depth 1 https://github.com/unknown-gd/quakejs-dockerized-ssl.git "${WORK_DIR}"; then
    bashio::exit.nok "git clone failed"
fi

cd "${WORK_DIR}" || bashio::exit.nok "Failed to enter work dir"

# Log what actually cloned so we can see the real file layout
bashio::log.info "Cloned files in ${WORK_DIR}:"
ls -la "${WORK_DIR}/"

# Find compose.yml
COMPOSE_FILE=$(find "${WORK_DIR}" -maxdepth 2 -name "compose.yml" | head -1 || true)
if [ -z "${COMPOSE_FILE}" ]; then
    bashio::log.error "compose.yml not found in cloned repo!"
    bashio::log.info "Full tree:"
    find "${WORK_DIR}" -maxdepth 3 | head -40
    bashio::exit.nok "compose.yml missing"
fi
bashio::log.info "Found compose file: ${COMPOSE_FILE}"

# --------------------------------------------------------------------------
# Patch compose.yml: replace ALL "./" volume paths with absolute host paths
# --------------------------------------------------------------------------
bashio::log.info "Patching compose.yml volume paths..."
bashio::log.info "  Replacing './' with '${HOST_WORK_DIR}/'"
sed -i "s|      - \./|      - ${HOST_WORK_DIR}/|g" "${COMPOSE_FILE}"
bashio::log.info "Patched volume lines:"
grep -n " - ${HOST_SHARE}" "${COMPOSE_FILE}" | head -20 || true

# Add game log volume mount to dedicated server so metrics server can read it
touch "${WORK_DIR}/game.log"
bashio::log.info "Patching dedicated server game log volume..."
sed -i "/dedicated-server.*server\.cfg:\/app\/base\/cpma\/server\.cfg/{
    a\      - ${HOST_WORK_DIR}/game.log:/app/game.log
}" "${COMPOSE_FILE}"
bashio::log.info "Dedicated server volumes:"
grep -A5 "q3js-dedicated-server" "${COMPOSE_FILE}" | grep " - " || true

# Remap port 443 -> 8443 (443 is taken by Nginx Proxy Manager)
bashio::log.info "Remapping port 443 -> 8443 in compose.yml..."
sed -i 's|"443:443"|"8443:443"|g; s|- 443:443|- 8443:443|g' "${COMPOSE_FILE}"
bashio::log.info "Port mappings after patch:"
grep -n "443" "${COMPOSE_FILE}" || true

# --------------------------------------------------------------------------
# Patch ALL Dockerfiles: insert sed before npm install to strip quakejs-files
# --------------------------------------------------------------------------
bashio::log.info "Patching Dockerfiles to remove quakejs-files..."
find "${WORK_DIR}/build" -name "Dockerfile" 2>/dev/null | while read -r DF; do
    bashio::log.info "Patching: ${DF}"
    awk '{
        if (/npm install/ && !/quakejs-files/) {
            print "RUN sed -i \x27/\"quakejs-files\"/d\x27 package.json 2>/dev/null || true"
        }
        print
    }' "${DF}" > "${DF}.patched" && mv "${DF}.patched" "${DF}"
    bashio::log.info "  patch lines added: $(grep -c 'quakejs-files' "${DF}")"
done

bashio::log.info "Pruning Docker build cache..."
docker builder prune -af 2>/dev/null || true

# --------------------------------------------------------------------------
# Generate self-signed certificate
# --------------------------------------------------------------------------
CERT_DIR="${WORK_DIR}/config/certificates"
mkdir -p "${CERT_DIR}"

if [ ! -f "${CERT_DIR}/q3js.crt" ]; then
    bashio::log.info "Generating self-signed certificate..."
    openssl req -x509 -nodes -days 3650 -newkey rsa:2048 \
        -keyout "${CERT_DIR}/q3js.key" \
        -out "${CERT_DIR}/q3js.crt" \
        -subj "/CN=${HOST_IP}" 2>/dev/null
fi

cat > "${CERT_DIR}/config.json" << CERTEOF
{
    "cert": "/app/certificates/q3js.crt",
    "key": "/app/certificates/q3js.key"
}
CERTEOF

# --------------------------------------------------------------------------
# Copy pak assets from /share/q3js/assets if present
# --------------------------------------------------------------------------
mkdir -p "${WORK_DIR}/assets"
if [ -d "/share/q3js/assets" ]; then
    bashio::log.info "Copying assets from /share/q3js/assets..."
    cp -rn /share/q3js/assets/. "${WORK_DIR}/assets/" 2>/dev/null || true
fi

# --------------------------------------------------------------------------
# Write server.cfg
# --------------------------------------------------------------------------
mkdir -p "${WORK_DIR}/config/dedicated-server/baseq3"
mkdir -p "${WORK_DIR}/config/dedicated-server/cpma"

cat > "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg" << CFGEOF
seta sv_hostname "HA QuakeJS Server"
seta sv_maxclients 12
seta g_motd "Home Assistant Quake 3"
seta g_gametype 0
seta fraglimit ${FRAG_LIMIT}
seta timelimit ${TIME_LIMIT}
seta rconpassword "${RCON_PASS}"
seta g_log "/app/game.log"
seta g_logSync 1
seta bot_enable 1
CFGEOF

BOT_COUNT=$(bashio::config 'bots | length')
for i in $(seq 0 $((BOT_COUNT - 1))); do
    BOT_NAME=$(bashio::config "bots[${i}].name")
    BOT_SKILL=$(bashio::config "bots[${i}].skill")
    echo "addbot ${BOT_NAME} ${BOT_SKILL}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"
done
echo "map ${MAP_NAME}" >> "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg"

cp "${WORK_DIR}/config/dedicated-server/baseq3/server.cfg" \
   "${WORK_DIR}/config/dedicated-server/cpma/server.cfg"

# --------------------------------------------------------------------------
# Write .env
# --------------------------------------------------------------------------
cat > "${WORK_DIR}/.env" << ENVEOF
Q3JS_DOMAIN=${HOST_IP}
Q3JS_SERVER_GAMEMODE=${GAMEMODE}
Q3JS_SERVER_DEDICATED=1
Q3JS_SERVER_RESTART_HOURS=${RESTART_HOURS}
ENVEOF

# --------------------------------------------------------------------------
# Log the final volume section for debugging
# --------------------------------------------------------------------------
bashio::log.info "Final volume mounts in compose.yml:"
grep -A1 "volumes:" "${COMPOSE_FILE}" | grep " - /" | head -20 || true

# --------------------------------------------------------------------------
# Build and start via docker compose
# --------------------------------------------------------------------------
bashio::log.info "Building Q3JS containers (first run takes several minutes)..."
docker compose -f "${COMPOSE_FILE}" build 2>&1 | tee "${LOG_DIR}/build.log"
if [ $? -ne 0 ]; then
    bashio::log.error "Build failed! Check /var/log/q3js/build.log"
    bashio::exit.nok "Build failed"
fi

bashio::log.info "Starting Q3JS containers..."
docker compose -f "${COMPOSE_FILE}" up -d 2>&1 | tee "${LOG_DIR}/compose.log"

bashio::log.info "Q3JS stack is up!"
bashio::log.info "  Web client : https://${HOST_IP}:8443"
bashio::log.info "  WSS proxy  : wss://${HOST_IP}:27961"
bashio::log.info "  Content    : https://${HOST_IP}:9000"

# --------------------------------------------------------------------------
# Start metrics server
# --------------------------------------------------------------------------
bashio::log.info "Starting metrics server on :8090..."
Q3_LOG_FILE="${WORK_DIR}/game.log" \
METRICS_PORT=8090 \
RCON_PASS="${RCON_PASS}" \
    /usr/bin/q3js-metrics > "${LOG_DIR}/metrics.log" 2>&1 &

# Keep alive
while true; do
    sleep 30
done
